# packaging/debian/postinst
#!/bin/sh
set -e

blue='\033[1;34m'
green='\033[1;32m'
yellow='\033[1;33m'
bold='\033[1;1m'
all_off='\033[0m'

ask_yn() {
    # $1 = câu hỏi, $2 = giá trị mặc định (y/n)
    local question="$1"
    local default="$2"
    local answer
    if [ "$default" = "y" ]; then
        printf "%b" "${bold}${question} [Y/n]: ${all_off}"
    else
        printf "%b" "${bold}${question} [y/N]: ${all_off}"
    fi
    # Nếu không phải terminal tương tác thì dùng mặc định
    if [ ! -t 0 ]; then
        echo "$default"
        answer="$default"
    else
        read -r answer
        answer=$(echo "$answer" | tr '[:upper:]' '[:lower:]')
        [ -z "$answer" ] && answer="$default"
    fi
    case "$answer" in
        y|yes) return 0 ;;
        *)     return 1 ;;
    esac
}

if [ "$1" = "configure" ]; then
    modprobe uinput || true

    # Reload udev
    if command -v udevadm >/dev/null; then
        udevadm control --reload-rules || true
        udevadm trigger || true
    fi
    systemd-sysusers || true

    echo ""
    echo -e "${blue}╔══════════════════════════════════════════╗${all_off}"
    echo -e "${blue}║       Thiết lập sau cài đặt Lotus        ║${all_off}"
    echo -e "${blue}╚══════════════════════════════════════════╝${all_off}"
    echo ""

    # --- Bước 1: Bật Lotus Server ---
    echo -e "${bold}[Bước 1/3] Bật Lotus Server${all_off}"
    echo -e "  Server giúp bộ gõ giao tiếp với hệ thống để xóa và sửa ký tự chính xác."
    echo -e "  ${yellow}(Yêu cầu quyền sudo)${all_off}"
    echo ""
    if ask_yn "  Bật Lotus Server cho user \"$(whoami)\" ngay bây giờ?" "y"; then
        echo ""
        if sudo systemctl enable --now "fcitx5-lotus-server@$(whoami).service" 2>/dev/null || \
           (sudo systemd-sysusers && sudo systemctl enable --now "fcitx5-lotus-server@$(whoami).service" 2>/dev/null); then
            echo -e "  ${green}✓ Lotus Server đã được bật thành công.${all_off}"
        else
            echo -e "  ${yellow}⚠ Không thể bật tự động. Chạy thủ công:${all_off}"
            echo -e "    ${blue}sudo systemctl enable --now fcitx5-lotus-server@\$(whoami).service${all_off}"
        fi
    else
        echo ""
        echo -e "  ${yellow}→ Bỏ qua. Chạy thủ công sau khi cài đặt:${all_off}"
        echo -e "    ${blue}sudo systemctl enable --now fcitx5-lotus-server@\$(whoami).service${all_off}"
    fi
    echo ""

    # --- Bước 2: Thiết lập biến môi trường ---
    echo -e "${bold}[Bước 2/3] Thiết lập biến môi trường${all_off}"
    echo -e "  Các biến GTK_IM_MODULE, QT_IM_MODULE, XMODIFIERS... cho phép ứng dụng"
    echo -e "  nhận diện Fcitx5 là bộ gõ mặc định. Thiếu các biến này bộ gõ sẽ không hoạt động."
    echo ""
    if ask_yn "  Thêm biến môi trường Fcitx5 vào \$HOME?" "y"; then
        echo ""
        # Phát hiện shell mặc định
        USER_SHELL=$(getent passwd "$(whoami)" | cut -d: -f7)
        case "$USER_SHELL" in
            */zsh)
                PROFILE_FILE="$HOME/.zprofile"
                ;;
            */fish)
                PROFILE_FILE="$HOME/.config/fish/config.fish"
                ;;
            *)
                PROFILE_FILE="$HOME/.bash_profile"
                ;;
        esac

        IM_VARS='GTK_IM_MODULE=fcitx
QT_IM_MODULE=fcitx
XMODIFIERS=@im=fcitx
SDL_IM_MODULE=fcitx
GLFW_IM_MODULE=ibus'

        # Kiểm tra xem đã có chưa để tránh thêm trùng
        if grep -q 'GTK_IM_MODULE=fcitx' "$PROFILE_FILE" 2>/dev/null; then
            echo -e "  ${yellow}⚠ Biến môi trường đã tồn tại trong ${PROFILE_FILE}, bỏ qua.${all_off}"
        else
            if [ "$USER_SHELL" = "$(command -v fish 2>/dev/null)" ] || echo "$USER_SHELL" | grep -q fish; then
                mkdir -p "$(dirname "$PROFILE_FILE")"
                cat >> "$PROFILE_FILE" <<'FISHEOF'
if status is-login
    set -Ux GTK_IM_MODULE fcitx
    set -Ux QT_IM_MODULE fcitx
    set -Ux XMODIFIERS "@im=fcitx"
    set -gx SDL_IM_MODULE fcitx
    set -gx GLFW_IM_MODULE ibus
end
FISHEOF
            else
                {
                    echo ""
                    echo "# Fcitx5 input method environment variables (added by fcitx5-lotus)"
                    echo "export GTK_IM_MODULE=fcitx"
                    echo "export QT_IM_MODULE=fcitx"
                    echo "export XMODIFIERS=@im=fcitx"
                    echo "export SDL_IM_MODULE=fcitx"
                    echo "export GLFW_IM_MODULE=ibus"
                } >> "$PROFILE_FILE"
            fi
            echo -e "  ${green}✓ Đã thêm vào ${PROFILE_FILE}.${all_off}"
            echo -e "  ${yellow}→ Log out và log in lại để áp dụng.${all_off}"
        fi
    else
        echo ""
        echo -e "  ${yellow}→ Bỏ qua. Xem README để thêm thủ công.${all_off}"
    fi
    echo ""

    # --- Bước 3: Tắt IBus ---
    echo -e "${bold}[Bước 3/3] Tắt bộ gõ cũ (IBus)${all_off}"
    echo -e "  IBus và Fcitx5 cùng chạy sẽ xung đột, gây gõ đôi ký tự hoặc bộ gõ không hoạt động."
    echo -e "  Bước này sẽ dừng ibus-daemon và vô hiệu hóa IBus khỏi autostart."
    echo ""
    if command -v ibus-daemon >/dev/null 2>&1 || command -v ibus >/dev/null 2>&1; then
        if ask_yn "  Tắt IBus ngay bây giờ?" "y"; then
            echo ""
            killall ibus-daemon 2>/dev/null || ibus exit 2>/dev/null || true
            # Xóa autostart IBus nếu có
            rm -f "$HOME/.config/autostart/ibus.desktop" \
                  "$HOME/.config/autostart/ibus-daemon.desktop" 2>/dev/null || true
            echo -e "  ${green}✓ IBus đã được tắt.${all_off}"
            echo -e "  ${yellow}→ Hãy thêm Fcitx5 vào autostart của DE/WM để khởi động cùng hệ thống.${all_off}"
        else
            echo ""
            echo -e "  ${yellow}→ Bỏ qua. Nhớ tắt IBus thủ công nếu gặp xung đột.${all_off}"
        fi
    else
        echo -e "  ${green}✓ Không phát hiện IBus, bỏ qua bước này.${all_off}"
    fi
    echo ""

    echo -e "${blue}─────────────────────────────────────────────${all_off}"
    echo -e "${bold}Các bước tiếp theo:${all_off}"
    echo -e "  4. Mở ${green}Fcitx5 Configuration${all_off} (fcitx5-configtool), thêm bộ gõ ${green}Lotus${all_off}."
    echo -e "  5. ${yellow}Wayland (KDE):${all_off} Vào System Settings → Keyboard → Virtual Keyboard → chọn Fcitx 5."
    echo -e "  6. ${yellow}Hyprland:${all_off} Thêm vào hyprland.conf:"
    echo -e "       ${blue}permission = fcitx5-lotus-server, keyboard, allow${all_off}"
    echo -e "${blue}─────────────────────────────────────────────${all_off}"
    echo ""
fi

#DEBHELPER#

exit 0
